# bot.conf
[program:ch_bot]
command=/home/ch_bot/venv/bin/gunicorn ch_bot.wsgi:application -c /home/ch_bot/ch_bot/gunicorn.conf.py
directory=/home/ch_bot
user=nobody
autorestart=true
redirect_stderr=true